# node-react/roles/deploy_node_app/files/node_app.service
[Unit]
Description=Node.js Application
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/app
ExecStart=/usr/bin/npm start
Restart=always

[Install]
WantedBy=multi-user.target