# Message_Board_MERN/labDirectory/roles/deploy_app/files/node_app.service
[Unit]
Description=Node.js Application
After=network.target mongod.service
Requires=mongod.service

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/app
ExecStart=/usr/bin/npm start
Restart=always

[Install]
WantedBy=multi-user.target